# upstart/taddws.conf
description "Torrent Add Web Service"

start on runlevel [2345]

respawn

exec ruby /usr/local/bin/taddws 
