# upstart/taddws.conf
description "Torrent Add Web Service"

start on started network
stop on stopped network
stop on starting shutdown

respawn

exec ruby /usr/local/bin/taddws 
